# puppet/modules/jupyter/templates/jupyter.conf.erb
[program:jupyter]
command = jupyter notebook --ip=* --no-browser
autostart = true
autorestart = true
user = <%= @user %>
directory = <%= @directory %>
